# devops/gunicorn_supervisor.ini
[program:gunicorn]
command=/usr/local/envs/myenv/bin/gunicorn --pythonpath /opt/apps/ --log-level debug app:app -c /opt/apps/devops/gunicorn.conf.py
directory=/opt/apps/
user=root
autostart=true
autorestart=true
redirect_stderr=true
environment =
    MYENV1="MYENV1",
    MYENV2="MYENV2"
